# samples/moon.service
[Unit]
Description=Moon - Dynamic Headless Engine
Documentation=https://github.com/thalib/moon
After=network.target

[Service]
Type=simple
User=moon
Group=moon
WorkingDirectory=/var/lib/moon

# Environment variables
# IMPORTANT: Set MOON_JWT_SECRET before starting the service
# You can set it by running: systemctl edit moon.service
# Or by editing this file and adding: Environment="MOON_JWT_SECRET=your-secret-here"
Environment="MOON_CONFIG=/etc/moon.conf"

# Execute the Moon binary
ExecStart=/usr/local/bin/moon

# Restart policy
Restart=on-failure
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=moon

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/moon
PrivateDevices=true

[Install]
WantedBy=multi-user.target
